# docker/node/supervisor/serve.ini
[program:serve]
process_name=%(program_name)s
command=/app/node_modules/.bin/vue-cli-service serve --port=9100 --host %(ENV_VUE_BIND_HOST)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:translation_observe]
process_name=%(program_name)s
command=node /app/bin/translations/translations-observer.js
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
