# docker/crontabs/schedule-cron
#!/bin/sh
* * * * * php /app/artisan schedule:run >> /dev/stdout 2>&1
